# crontab
0 * * * * $HOME/bin/change_wallpapers.sh
